# edge/raspberry_pi/system/edge-ai.service
[Unit]
Description=Edge AI Irrigation Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/edgeai/app
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /home/pi/edgeai/app/edge_ai_service.py
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target